# packages/sysutils/mdadm/udev/rules.d/md-mount.rules
# Automatically mount Linux RAID MD devices
ACTION=="add", SUBSYSTEM=="block", KERNEL=="md*", ENV{DEVTYPE}=="disk", RUN+="/usr/bin/mount-md-device.sh %k"
